# VMs/frontend/website/emailCronJob.crontab
0 0 1 * * /usr/bin/node /IT490-Project/VMs/frontend/website/emailCronJob.js
