# VMs/backend/services/reviews-receiver/reviews-receiver.service
[Unit]
Description=RabbitMQ Search Receiver Service for Reviews
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/php /opt/reviews-receiver/RabbitMQServer.php
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
