# VMs/database/services/discussion-receiver/discussion-receiver.service
[Unit]
Description=RabbitMQ Search Receiver Service for Discussions
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/php /opt/discussion-receiver/RabbitMQServer.php
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target